# _linux/var/spool/cron/crontabs/root
* * * * * /usr/bin/php /wwwroot/cjkbe/artisan schedule:run >> /dev/null 2>&1
